# yapo-scraper-telegram-bot.timer
[Unit]
Description=Run yapo scraper each 10 min

[Timer]
OnUnitActiveSec=5
OnBootSec=5s

[Install]
WantedBy=timers.target
